# keycloak_addon/rootfs/etc/services.d/keycloak_addon/run
#!/usr/bin/with-contenv bashio

set -e
bashio::log.info "== Keycloak Add-on starting =="

# Wait for MariaDB
bashio::log.info "Looking for MariaDB..."
if ! bashio::services.available 'mysql'; then
    bashio::log.fatal "Local database access should be provided by the MariaDB addon"
    bashio::exit.nok "Please ensure it is installed and started"
fi

DB_ADDR=$(bashio::services "mysql" "host")
DB_PORT=$(bashio::services "mysql" "port")
DB_DATABASE=$(bashio::config 'db_name')
DB_USER=$(bashio::config 'db_user')
DB_PASSWORD=$(bashio::config 'db_password')
export KEYCLOAK_ADMIN=$(bashio::config 'KEYCLOAK_ADMIN')
export KEYCLOAK_ADMIN_PASSWORD=$(bashio::config 'KEYCLOAK_ADMIN_PASSWORD')

# Вычисляем HTTP relative path для Ingress
REL_PATH=$(bashio::addon.ingress_entry)
bashio::log.info "Using relative path: ${REL_PATH}"

## Run your program
exec /opt/keycloak/bin/kc.sh start-dev \
  --proxy-headers=xforwarded \
  --hostname-strict=false \
  --http-relative-path="${REL_PATH}" \
  --db=mariadb \
  --db-url=jdbc:mariadb://"${DB_ADDR}":"${DB_PORT}"/"${DB_DATABASE}" \
  --db-username="${DB_USER}" \
  --db-password="${DB_PASSWORD}" \
  --features=token-exchange
