# vikunja_addon/rootfs/etc/services.d/vikunja_addon/finish
#!/usr/bin/env bashio
bashio::log.info "Vikunja service stopped"
exit 0